# deploy/systemd/garagebrain-api.service
[Unit]
Description=GarageBrain API
After=network.target postgresql.service garagebrain-gateway.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/garagebrain
ExecStart=/var/www/garagebrain/bin/server
Restart=always
RestartSec=5
EnvironmentFile=/var/www/garagebrain/.env
MemoryLimit=200M

[Install]
WantedBy=multi-user.target
